# systemd/qs-security-scan.service
[Unit]
Description=Bar security-status scan (AUR-Malware + bumblebee)

[Service]
Type=oneshot
ExecStart=%h/.local/bin/qs-security-scan.sh
Nice=19
IOSchedulingClass=idle
