# systemd/qs-shell-update-check.timer
[Unit]
Description=QS-Shell update check timer

[Timer]
OnBootSec=2min
# also fire shortly after the timer itself is (re)started — without this a
# daemon-reload/re-enable mid-session leaves the timer "elapsed" with no NEXT
OnActiveSec=2min
OnUnitActiveSec=6h
Persistent=true
RandomizedDelaySec=10min

[Install]
WantedBy=timers.target
